# nextbus/crontab
SHELL=/bin/sh
PATH=/usr/local/sbin:/usr/local/bin:/sbin:/bin:/usr/sbin:/usr/bin

* * * * * /src/nextbus.sh >> /var/log/cron.log 2>&1
#
